# debian/postinst
#!/bin/bash

set -e

. /etc/default/openmediavault
. /usr/share/openmediavault/scripts/helper-functions

case "$1" in
  configure)
    # Activate package triggers.
    dpkg-trigger update-workbench

    # Initialize configuration database.
    omv-confdbadm create "conf.service.bcachefs"
    if [ -n "$2" ]; then
      omv-confdbadm migrate "conf.service.bcachefs" "${2}"
    fi

    # Enable the auto-unlock service (no-op if no key files exist yet).
    systemctl enable bcachefs-autounlock.service || true

    # Add the bcachefs apt repository
    omv-salt deploy run --no-color apt || true

    # Run bcachefs salt states (package install + cron setup).
    omv-salt deploy run --no-color bcachefs bcachesfcron || true

    # Mark the bcachefs module dirty if the kernel module is not loaded,
    # so the user is prompted to apply changes to install bcachefs-tools
    # (and bcachefs-kernel-dkms if needed).
    if ! lsmod | grep -q "^bcachefs "; then
      omv_module_set_dirty bcachefs
    fi
  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
  ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

#DEBHELPER#

exit 0
